# demo/deploy/templates/postgresql.ini


[program:postgresql]
command=postmaster -D "{{config_prefix}}/etc/"
process_name=postgresql
stdout_logfile = {{config_logdir}}/postgresql.log
stderr_logfile = {{config_logdir}}/postgresql.err
pidfile = {{config_prefix}}/var/run/postgresql.pid
stopsignal=INT
autostart=true
autorestart=true
redirect_stderr=true